# supervisor.ini
[program:ocserv-443]
command= ocserv -c /etc/ocserv/ocserv.conf -f
autostart = true
startsecs = 5
autorestart = true
startretries = 30
user = root
redirect_stderr = true
stdout_logfile = /tmp/ocserv-443.log

[program:sss-443]
command= ss-redir -c /etc/shadowsocks.json -u -b 0.0.0.0
autostart = true
startsecs = 5
autorestart = true
startretries = 30
user = root
redirect_stderr = true
stdout_logfile = /tmp/ocserv-443.log

[program:ocserv-443]
command= dnsmasq -k
autostart = true
startsecs = 5
autorestart = true
startretries = 30
user = root
redirect_stderr = true
stdout_logfile = /tmp/dnsmasq-1053.log

